# crontab
# clean up links created by my OpenWith script
0 0-23 * * * find public_html/tmp/ -type l -mtime +1 -exec rm "{}" \;

